# roles/consul/files/backup.cron
PATH=/usr/local/sbin:/sbin:/usr/local/bin:/usr/bin:/bin
30 0 * * 0 root consul_clean
30 0 * * * root consul_backup

