# roles/consul/files/backup.cron
30 0 * * sun root consul_clean
30 0 * * * root consul_backup

